# docs/research/routellm-phase3/classifier/launchd/com.nirarin.router-drift.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nirarin.router-drift</string>
    <key>ProgramArguments</key>
    <array>
        <string>__MANIFESTO_ROOT__/scripts/router-drift-check.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__MANIFESTO_ROOT__/docs/research/routellm-phase3/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__MANIFESTO_ROOT__/docs/research/routellm-phase3/logs/launchd-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
